# 1.0/src/app.sv.conf
[program:app]
command=/usr/local/bin/start_app.sh
autostart=true
autorestart=true
